# service/systemd/morbled.service
[Unit]
Description=Screen backlight daemon
Requires=dbus

[Service]
Type=exec
ExecStart=morbled
